# alpine/agent
#!/sbin/openrc-run

command="/cocos/agent"
command_args=""

depend() {
    need net
}

start_pre() {
    mkdir -p /run/cocos
}

start() {
    AGENT_GRPC_ADDR=$(ifconfig eth0 | grep "inet addr" | cut -d ':' -f 2 | cut -d ' ' -f 1):7002
    ebegin "Starting cocos agent"
    env AGENT_GRPC_ADDR="${AGENT_GRPC_ADDR}" start-stop-daemon -S -q -b -m -p /run/cocos/agent.pid -x $command -- $command_args
    eend $?
}

stop() {
    ebegin "Stopping cocos agent"
    start-stop-daemon -K -q -p /run/cocos/agent.pid
    eend $?
}
