# init/systemd/cocos-manager.service
[Unit]
Description=Cocos Manager
After=network.target

[Service]
ExecStart=cocos-manager
Restart=always
RestartSec=5s
EnvironmentFile=/etc/cocos/cocos-manager.env

[Install]
WantedBy=multi-user.target
