# init/systemd/cocos-agent.service
[Unit]
Description=Cocos AI agent
After=network.target attestation-service.service
Requires=attestation-service.service
Before=docker.service

[Service]
WorkingDirectory=/cocos
StandardOutput=file:/var/log/cocos/agent.stdout
StandardError=file:/var/log/cocos/agent.stderr
EnvironmentFile=/etc/cocos/environment
ExecStartPre=/cocos_init/agent_setup.sh
ExecStart=/cocos_init/agent_start_script.sh
Restart=always
RestartSec=5s

[Install]
WantedBy=default.target
