# hal/disk/board/rootfs-overlay/usr/lib/systemd/system/cocos-agent.service.d/mounts.conf
[Unit]
RequiresMountsFor=/etc/cocos /etc/certs

